# metrics-server/huginn-metrics.service
[Unit]
Description=Huginn Metrics Collector & API
After=network.target nginx.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/huginn-metrics
ExecStart=/usr/bin/node src/index.js
Restart=always
RestartSec=5

Environment=PORT=3100
Environment=NGINX_LOG_PATH=/var/log/nginx/huginn_metrics.log
Environment=API_KEY=
Environment=ALLOWED_ORIGINS=https://status-monad-huginn.tech,https://status-monad-testnet-huginn.tech

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=huginn-metrics

[Install]
WantedBy=multi-user.target
